# halt.service
[Unit]
Description=Halt systemd

[Service]
ExecStart=systemctl --force halt
